# mqtt-s7-connector/rootfs/etc/s6-overlay/s6-rc.d/config-gui/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: MQTT S7 Connector
# Startet die Konfigurations-GUI
# ==============================================================================

CONFIG_GUI_FILE="/addon_configs/$(bashio::addon.slug)/config.yaml"
export CONFIG_GUI_FILE
export CONFIG_GUI_TEMPLATE="/usr/lib/mqtt-s7-config-gui/config.template.yaml"

bashio::log.info "Starte Konfigurations-GUI auf Port ${CONFIG_GUI_PORT:-8099}"

exec node /usr/lib/mqtt-s7-config-gui/server.js
